# test
#!/bin/sh /etc/rc.common

#just a start no stop
START =100
start() {
    lua /root/ubusPingService.lua &

}
